# docker/supervisord.conf
[supervisord]
nodaemon=true
user=ftuser
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
loglevel=info

[program:freqtrade]
command=/home/ftuser/.local/bin/freqtrade trade --config /freqtrade/user_data/config.prod.json --logfile /freqtrade/user_data/logs/freqtrade.log
directory=/freqtrade
user=ftuser
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=30
stdout_logfile=/var/log/supervisor/freqtrade.log
stderr_logfile=/var/log/supervisor/freqtrade_err.log
environment=PYTHONUNBUFFERED=1

[program:demo_ui]
command=/home/ftuser/.local/bin/python -m freqtrade.ui.demo_server
directory=/freqtrade
user=ftuser
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=10
stdout_logfile=/var/log/supervisor/demo_ui.log
stderr_logfile=/var/log/supervisor/demo_ui_err.log
environment=PYTHONUNBUFFERED=1,FLASK_RUN_HOST=0.0.0.0,FLASK_RUN_PORT=5000

[program:config_dashboard]
command=/home/ftuser/.local/bin/streamlit run freqtrade/ui/prod_config.py --server.port 8501 --server.address 0.0.0.0 --server.headless true
directory=/freqtrade
user=ftuser
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=10
stdout_logfile=/var/log/supervisor/config_dashboard.log
stderr_logfile=/var/log/supervisor/config_dashboard_err.log
environment=PYTHONUNBUFFERED=1

[program:monitor_dashboard]
command=/home/ftuser/.local/bin/streamlit run freqtrade/ui/prod_monitor.py --server.port 8502 --server.address 0.0.0.0 --server.headless true
directory=/freqtrade
user=ftuser
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=10
stdout_logfile=/var/log/supervisor/monitor_dashboard.log
stderr_logfile=/var/log/supervisor/monitor_dashboard_err.log
environment=PYTHONUNBUFFERED=1
